# supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/bin/sh -c "/usr/sbin/nginx"



